# deploy/staragent-ingest.timer
# Schedule for StarAgent knowledge-base ingestion.
# Default: daily at 04:17 local (off-peak; incremental thanks to hash-skip).
# Change OnCalendar to suit — e.g. "weekly", "Mon *-*-* 05:00:00", "*-*-* */6:00".

[Unit]
Description=Run StarAgent ingestion on a schedule

[Timer]
OnCalendar=__ONCALENDAR__
# Spread load / avoid the exact minute.
RandomizedDelaySec=600
# If the machine was off at the scheduled time, run once at next boot.
Persistent=true
Unit=staragent-ingest.service

[Install]
WantedBy=timers.target
